# files/etc/crontabs/root
0 6 * * * /usr/share/adbyby/adblock.sh > /tmp/adupdate.log 2>&1
